# etc/skel/.config/systemd/user/multi-user.target.wants/xdg_menu.path
[Unit]
Description=Monitora arquivos modificados na pasta especificada
​
[Path]
PathChanged=/usr/share/applications/
​
[Install]
WantedBy=multi-user.target
